# install-scripts/components/gnosisvpn/seqs-gnosisvpn-dns.path
[Unit]
Description=Watch for GnosisVPN resolver changes

[Path]
PathChanged=/run/resolvconf/interfaces
Unit=seqs-gnosisvpn-dns.service

[Install]
WantedBy=multi-user.target
